# ansible/blue-team-tools/mozdef/roles/mozdef-install/files/mozdef.service
[Unit]
Description=mozdef service

[Service]
WorkingDirectory=/opt/mozdef
ExecStart=make run

[Install]
WantedBy=multi-user.target
