# ansible/graylog/roles/firewall-config/files/iptables-restore.service
[Unit]
Description=Iptables service

[Service]
User=root
ExecStart=/opt/iptables-restore.sh

[Install]
WantedBy=multi-user.target
